# integrations/ostree/bdfs-ostree-prune.service
# systemd service — prune old OSTree deployments after each OS update.
# Install to /etc/systemd/system/ and enable with:
#   systemctl enable --now bdfs-ostree-prune.timer

[Unit]
Description=Prune old OSTree deployments (bdfs-ostree)
After=ostree-finalize-staged.service
Wants=ostree-finalize-staged.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bdfs-ostree prune
# Source site config if present
EnvironmentFile=-/etc/bdfs/ostree.conf

[Install]
WantedBy=multi-user.target
